# tests/apps/app1/systemd/app@.target
[Unit]
Description={{ settings.PROJECT_NAME|default:"Django" }} Instance %i (app1)

[Install]
WantedBy=multi-user.target
